# upstream/distributed-ensemble/java-intellij/idea.conf
[program:idea]
command=/home/ensemble/idea/bin/idea.sh
autorestart=false
